# scripts/systemd/freddy-backup.timer
[Unit]
Description=FREDDY Backup Timer
Documentation=https://github.com/nuniesmith/personal_servers_freddy
Requires=freddy-backup.service

[Timer]
# Run daily at 2 AM
OnCalendar=daily
OnCalendar=*-*-* 02:00:00

# Run 10 minutes after boot if missed
Persistent=true

# Randomize start time by up to 30 minutes
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
